# assets/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:php7-fpm]
command=/usr/sbin/php-fpm7.0
autorestart=true

[program:cron]
command=cron -f

[program:sshd]
command=/usr/sbin/sshd
autorestart=true